# roles/docker-install/files/docker.service
[Unit]
Description=docker

[Service]

ExecStart=/usr/bin/dockerd
PIDFile=/var/run/dockerd

TimeoutStartSec=600
TimeoutStopSec=600
